# supervisor.conf.d/redsocks.conf
[program:redsocks]
command=redsocks -c /etc/redsocks.conf
