# supervisor.conf.d/squid.conf
[program:squid]
command=squid3 -N
startsecs=10
stopwaitsecs=10
